# deploy/claude-automation.service
[Unit]
Description=Claude Automation Service
After=network.target

[Service]
User=claude
Group=claude
WorkingDirectory=/opt/claude-automation
ExecStart=/usr/local/bin/bun run /opt/claude-automation/src/server.ts
Restart=always
Environment="BUN_ENV=production"
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target
